# binder.service
[Unit]
Description=Discord Bot
After=network.target

[Service]
Type=simple
ExecStart=$BINARY
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
